# docker-wait
#!/sbin/openrc-run
command="/usr/bin/true"
description="Wait for network to get IP and then start docker"

depend() {
    need net
    after wpa_supplicant NetworkManager  # опционально, можно убрать или добавить нужный менеджер
}

start() {
    ebegin "Waiting for network to be up (max 60s)"
    for i in $(seq 1 60); do
        # проверяем есть ли IPv4 адрес на интерфейсах
        if ip -4 addr show scope global | grep -q "inet "; then
            eend 0
            ebegin "Starting docker"
            doas rc-service docker start
            eend $?
            return 0
        fi
        sleep 1
    done
    eend 1
    return 1
}
